# BLECarPiZeroW/scripts/picamera_service/picamera.service
[Unit]
Description=Video Streaming service using PiCamera
After=network.target

[Service]
User=root
Group=pidevs
IgnoreSIGPIPE=false
ExecStart=/usr/bin/picamera_start

[Install]
WantedBy=multi-user.target
